# crontab
* * * * * php /mnt/app/index.php

# always remeber the last line of a crontab file must eb an empty line

